# components/common/salt/state/monitoring/lit-os-metrics-host/etc/systemd/system/lit-os-metrics-host-15min.timer.j2
[Unit]
Description=Run Lit OS Metrics Host every 15 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=15min
AccuracySec=1s

[Install]
WantedBy=timers.target
