# components/common/salt/state/monitoring/lit-os-metrics-host/etc/systemd/system/lit-os-metrics-host-1min.timer.j2
[Unit]
Description=Run Lit OS Metrics Host every 1 minute

[Timer]
OnBootSec=60
OnUnitActiveSec=1min
AccuracySec=1s

[Install]
WantedBy=timers.target
